# deploy/systemd/nami-worker@.service
[Unit]
Description=Nami Worker — %i
After=network.target nami-core.service
PartOf=nami-core.service

[Service]
Type=simple
User=nami-worker
WorkingDirectory=/opt/nami-core
Environment=PYTHONPATH=/opt/nami-core/src
Environment=NAMI_CONFIG_DIR=/opt/nami-core/config
Environment=NAMI_WORKER=%i
EnvironmentFile=/etc/nami-harness/nami-core.env
ExecStart=/opt/nami-core/.venv/bin/python -m nami_core.runtime.queue.worker
Restart=always
RestartSec=5
WatchdogSec=60
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
